# voice_mode/templates/systemd/voicemode-kokoro.service
# voicemode-kokoro.service v1.3.0
# Last updated: 2026-02-18
# Compatible with: kokoro-fastapi v1.0.0+
# Simplified: start script handles config via voicemode.env

[Unit]
Description=Voice Mode Kokoro Text-to-Speech Service
After=network.target

[Service]
Type=simple
ExecStart={START_SCRIPT}
WorkingDirectory={KOKORO_DIR}
Restart=on-failure
RestartSec=10
# Don't restart if the executable is missing
RestartPreventExitStatus=127

# Environment - start script sources voicemode.env for port config
# Include .cargo/bin for ARM64 Linux where Rust may be needed to build dependencies
Environment="PATH=%h/.cargo/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin"
# Restart worker after N requests to mitigate memory leak (hexgrad/kokoro#152)
Environment="UVICORN_LIMIT_MAX_REQUESTS={KOKORO_MAX_REQUESTS}"

# Resource limits
MemoryLimit=4G
CPUQuota=100%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voicemode-kokoro

[Install]
WantedBy=default.target