# voice_mode/templates/systemd/voicemode-whisper.service
# voicemode-whisper.service v1.1.0
# Last updated: 2025-11-12
# Uses unified startup script for dynamic model selection

[Unit]
Description=Whisper.cpp Speech Recognition Server
After=network.target

[Service]
Type=simple
ExecStart={START_SCRIPT_PATH}
Restart=on-failure
RestartSec=10
WorkingDirectory={INSTALL_DIR}
StandardOutput=append:{LOG_DIR}/whisper/whisper.out.log
StandardError=append:{LOG_DIR}/whisper/whisper.err.log
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/usr/local/cuda/bin"

[Install]
WantedBy=default.target