# voice_mode/templates/systemd/voicemode-whisper.service
# voicemode-whisper.service v1.2.1
# Last updated: 2026-05-28
# Simplified: start script handles config via voicemode.env

[Unit]
Description=Whisper.cpp Speech Recognition Server
After=network.target

[Service]
Type=simple
ExecStart={START_SCRIPT}
# Restart=always (not on-failure) so a clean exit (e.g. a future request-limit
# workaround mirroring kokoro's UVICORN_LIMIT_MAX_REQUESTS) still triggers a
# restart. Defensive alignment with kokoro fix (GH-448, VM-1398).
Restart=always
RestartSec=10

# Environment - start script sources voicemode.env for port/model config
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/local/cuda/bin"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voicemode-whisper

[Install]
WantedBy=default.target